# build/gppmd-v1.0.1-dev.1_amd64/lib/systemd/system/gppmd.service
[Unit]
Description=GPU Power and Performance Manager Daemon

[Service]
ExecStart=/usr/bin/gppmd
Restart=always
User=root
Group=root
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=gppmd

[Install]
WantedBy=multi-user.target
